# backend/services/refmap_optimized_trajectories.service
[Unit]
Description=RefMap Optimized Trajectories Service
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=/home/refmap-imafusa/Documents/RefMap/backend/optimized_trajectories
ExecStart=/home/refmap-imafusa/miniconda3/bin/python app.py
Restart=on-failure
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
