# system_services/refmap_emissions.service
[Unit]
Description=RefMap Emissions Service
After=network.target

[Service]
Type=simple
User=george
WorkingDirectory=/home/george/Documents/Microlab/RefMap/-RefMap-WP5_Dashboard/backend/emissions
ExecStart=/tools/miniconda3/envs/RefMapDashboard/bin/python app.py
Restart=on-failure
Environment=FLASK_ENV=production
Environment=PORT=4005

[Install]
WantedBy=multi-user.target
